# services/mongodb_service
#!/sbin/openrc-run

depend(){
    need net
}
command="mongod"
command_args="-f /server/mongo_config.yml"
command_user="root:root"
pidfile="/run/mongodb/pid"

