# services/mongodb_service
#!/sbin/openrc-run

depend(){
    need net
}
name="mongodb-sevice"
command="mongod"
command_args="-f /server/mongo_config.yml"

