# services/redis_service
#!/sbin/openrc-run

depend(){
    need net
}
command="usr/bin/redis-server"
command_user="root:root"
pidfile="/run/redis/pid"
