# services/redis_service
#!/sbin/openrc-run

depend(){
    need net
}
name="redis-service"
command="usr/bin/redis-server"
